# docker/cron.d/test-last-job
0 * * * * %USERNAME% /usr/bin/php /var/www/pastell/bin/console app:daemon:notify-check >> /data/log/pastell-daemon-check.log 2>&1
