# ci-resources/cron.d/logrotate
*   *   *   *   *   root    /usr/sbin/logrotate -f /etc/logrotate.conf >> /data/log/logrotate.log 2>&1
