# ci-resources/cron.d/maj-libersign
10  1   *   *   *   root  php /var/www/pastell/batch/maj-libersign.php  >> /data/log/maj-libersign.log 2>&1
