# ci-resources/cron.d/optimize
0 23 * * * www-data /usr/bin/php /var/www/pastell/batch/optimize.php >> /data/log/optimize.log 2>&1
