# ci-resources/cron.d/testcron
* * * * * root /bin/date > /tmp/last-cron
