# ci-resources/cron.d/testcron
* * * * * %USERNAME% /bin/date > /data/log/last-cron
