# ci-resources/cron.d/test-last-job
0 * * * * %USERNAME% /usr/bin/php /var/www/pastell/installation/test-last-job.php >> /data/log/test-last-job.log 2>&1
