# ci-resources/cron.d/validca
10 10 * * * root curl -s https://validca.libriciel.fr/retrieve-validca.sh | bash -s /data/certificate > /dev/null 2>&1 && supervisorctl apache2 reload >/dev/null 2>&1
