# ci-resources/cron.d/validca
10  *   *   *   *   root  /usr/local/bin/recup_crl.sh /etc/s2low/ssl/ > /data/log/recup_crl.log 2>&1
