# ci-resources/cron.d/vider-journal-to-historique
0 20 * * * www-data /usr/bin/php /var/www/pastell/installation/vider-journal-to-historique.php >> /data/log/vider-journal.log 2>&1
